# app-misc/interception-tools/files/interception-tools.init
#!/sbin/openrc-run

depend() {
	want udev-settle
	after udev-settle
}

command=/usr/bin/udevmon
# command_args='-c /etc/interception/udevmon.yml'
command_args=''
start_stop_daemon_args='-N -20'
pidfile=/var/run/udevmon.pid
command_background=true
output_log='/var/log/udevmon.log'
error_log='/var/log/udevmon.err'

